# extrarootfs/etc/sv/busybox-syslogd/run
#!/bin/sh
exec 2>&1
exec busybox syslogd -n -O - > /dev/ttyAMA0 2>&1
